# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for DOMAIN

[Service]
Restart=on-failure
User=cheena
WorkingDirectory=/home/cheena/sites/DOMAIN
EnvironmentFile=/home/cheena/sites/DOMAIN/.env

ExecStart=/home/cheena/sites/DOMAIN/venv/bin/gunicorn \
    --bind unix:/tmp/DOMAIN.socket \
    superlists.wsgi:application

[Install]
WantedBy=multi-user.target